# recipes-support/edt/files/90-edt.rules
ACTION=="add", ATTR{device}=="0x00a3", RUN+="/usr/sbin/edtmknod", 
