# recipes-support/epix/files/90-pixci.rules
ACTION=="add", ATTR{device}=="0xe004", RUN+="/usr/sbin/pixcimknod", 
